# systemd/claude-rc-reconcile.timer
[Unit]
Description=Periodically reconcile Claude remote-control sessions (keep them alive)

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
